# target/udev/99-md-raid-device-lvm.rules
SUBSYSTEM!="block", GOTO="mdraid_device_lvm_end"
ENV{UDISKS_MD_MEMBER_NAME}=="*:nvmevol*", GOTO="mdraid_device_lvm_add"
GOTO="mdraid_device_lvm_end"

LABEL="mdraid_device_lvm_add"
ACTION=="add", SUBSYSTEM=="block", ENV{ID_FS_USAGE}=="raid", ENV{ID_FS_TYPE}=="linux_raid_member", RUN+="/lib/lvm/lvm_mdraid.sh $devnode $env{UDISKS_MD_MEMBER_NAME}"

LABEL="mdraid_device_ltgt_end"
